# ubuntu-18.04/etc/systems/system/gridlabd.service
[Unit]
Description=GridLAB-D realtime simulator
After=network.target

[Service]
Type=simple
Restart=always
ExecStart=/usr/local/bin/gridlabd-rt

[Install]
WantedBy=multi-user.target
